# collablab-log.timer
[Unit]
Description=Run collablab-log every minute

[Timer]
OnCalendar=minutely
Unit=collablab-log.service

[Install]
WantedBy=timers.target
